# test/fixtures/unit-0.service
[Unit]
Description=Sleep infinity service
[Service]
ExecStart=sleep infinity
[Install]
WantedBy=default.target
